# n8n_onepager/deploy/systemd-n8nai-waitlist.service.example
[Unit]
Description=n8nAI Waitlist Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/n8nai.io/server
Environment=NODE_ENV=production
EnvironmentFile=/var/www/n8nai.io/server/.env
ExecStart=/usr/bin/node /var/www/n8nai.io/server/index.js
Restart=always
User=n8nai
Group=www-data

[Install]
WantedBy=multi-user.target

